# contrib/init/titancoind.conf
description "Titancoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env TTND_BIN="/usr/bin/titancoind"
env TTND_USER="titancoin"
env TTND_GROUP="titancoin"
env TTND_PIDDIR="/var/run/titancoind"
# upstart can't handle variables constructed with other variables
env TTND_PIDFILE="/var/run/titancoind/titancoind.pid"
env TTND_CONFIGFILE="/etc/titancoin/titancoin.conf"
env TTND_DATADIR="/var/lib/titancoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # titancoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$TTND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run titancoind."
        echo "The setting must appear in $TTND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Titancoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$TTND_PIDDIR"
    chmod 0755 "$TTND_PIDDIR"
    chown $TTND_USER:$TTND_GROUP "$TTND_PIDDIR"
    chown $TTND_USER:$TTND_GROUP "$TTND_CONFIGFILE"
    chmod 0660 "$TTND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$TTND_PIDFILE" \
    --chuid $TTND_USER:$TTND_GROUP \
    --exec "$TTND_BIN" \
    -- \
    -pid="$TTND_PIDFILE" \
    -conf="$TTND_CONFIGFILE" \
    -datadir="$TTND_DATADIR" \
    -disablewallet \
    -daemon

